# files/etc/systemd/system/bluetooth.service.d/override-a2dp.conf
[Unit]
ConditionPathIsDirectory=/sys/class/bluetooth/hci0
ConditionFileNotEmpty=/sys/class/bluetooth/hci0/name

[Service]
WatchdogSec=10
Restart=on-failure
